# rag-bi-platform.service
[Unit]
Description=RAG BI Platform Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/RAG-based-BI-Platform
ExecStart=/path/to/RAG-based-BI-Platform/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --threads 4 --timeout 120 run:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
